# embedly/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn embedly:app -b localhost:7001
directory=/embedly
